# deploy/systemd/event-crm-sync.timer
[Unit]
Description=Run Event CRM Hermes sync every minute

[Timer]
OnBootSec=30s
OnUnitActiveSec=60s
AccuracySec=5s
Unit=event-crm-sync.service

[Install]
WantedBy=timers.target

